# ik-borg-backup.service
[Unit]
Description=Isak Borg Backup (daily)

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ik-borg-backup
